# systemd/synchronizzer.service.erb
[Unit]
Description=Synchronizzer

[Service]
SyslogIdentifier=synchronizzer
WorkingDirectory=<%=@package_dir%>
Environment=DOTNET_CLI_TELEMETRY_OPTOUT=1
ExecStart=/usr/bin/dotnet <%=@package_dir%>/synchronizzer.dll --config <%=@config_file%> --hostingconfig <%=@hosting_file%>
ExecReload=/bin/kill -INT $MAINPID
Nice=10
ProtectSystem=full
ProtectHome=true
User=nobody
NoNewPrivileges=true
Restart=on-failure
RestartSec=47
